# ws.timer
[Unit]
Description=timer script for weather

[Timer]
OnCalendar=*-*-* *:00:00
Persistent=true

[Install]
WantedBy=timers.target
